# config/systemd/vector-search-monitor.service
[Unit]
Description=VANA Vector Search Monitoring Service
After=network.target

[Service]
User=nick
WorkingDirectory=/Users/nick/Development/vana
ExecStart=/usr/bin/python3 /Users/nick/Development/vana/scripts/monitor_vector_search.py
Restart=on-failure
StandardOutput=journal
StandardError=journal
EnvironmentFile=/Users/nick/Development/vana/.env

[Install]
WantedBy=multi-user.target