# config/systemd/vector-search-ui.service
[Unit]
Description=VANA Vector Search Dashboard UI Service
After=network.target vector-search-dashboard.service

[Service]
User=nick
WorkingDirectory=/Users/nick/Development/vana/dashboard
ExecStart=/Users/nick/Development/vana/venv/bin/streamlit run streamlit_app.py --server.port 8501 --server.address 0.0.0.0
Restart=on-failure
StandardOutput=journal
StandardError=journal
EnvironmentFile=/Users/nick/Development/vana/.env

[Install]
WantedBy=multi-user.target
